# monitor.service
[Unit]
Description=Monitor 'test' process and send heartbeat

[Service]
Type=oneshot
ExecStart=/home/ilzar/.local/monitor.sh
